# caddy.service
[service]
PrivateTmp=yes
Environment=CADDY_ADMIN=unix//tmp/caddy-admin.socket

[Install]
Requires=caddyplugins.socket