# cowsay.socket
[Unit]
Description=cowsay

[Socket]
ListenStream=9876
Accept=yes

[Install]
WantedBy=sockets.target
