# cowsay@.service
[Unit]
Description=cowsay

[Service]
ExecStart=/bin/cowsay
StandardInput=socket
